# overlay_fs/etc/inittab
# Ceci est un script pour init
::sysinit:/etc/init.d/rcS
ttyS0::askfirst:-/bin/sh
::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
